# scripts/wifi-fallback.service
[Unit]
Description=WiFi Fallback Hotspot Service
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/wifi_fallback.sh
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
